# aria2/root/etc/services.d/02-caddy/run
#!/usr/bin/execlineb -P

exec /usr/local/bin/caddy -quic -conf /etc/conf/caddy/Caddyfile
